# systemd/tailscale-online.target
[Unit]
Description=Tailscale is Online
Wants=tailscale-wait-online.service
After=tailscale-wait-online.service
