# systemd/tailscale-online.target
[Unit]
Description=Tailscale is Online
After=tailscale-wait-online.service
Wants=tailscale-wait-online.service
